# deploy/com.wordsbot.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.wordsbot.backup</string>
  <key>ProgramArguments</key>
  <array><string>/Users/misha/work/words-trainer/scripts/backup_postgres.sh</string></array>
  <key>WorkingDirectory</key><string>/Users/misha/work/words-trainer</string>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>15</integer></dict>
  <key>StandardOutPath</key><string>/Users/misha/Library/Logs/wordsbot/backup.out.log</string>
  <key>StandardErrorPath</key><string>/Users/misha/Library/Logs/wordsbot/backup.err.log</string>
</dict>
</plist>
